# contrib/deb/DEBIAN-BIN/prerm
#!/bin/sh
set -e
systemctl stop padctl.service padctl-resume.service || true
systemctl disable padctl.service padctl-resume.service || true
